# etc/udev/rules.d/99-xmodmap.rules
ACTION=="add", ATTR{idVendor}=="04fe", ATTR{idProduct}=="0022", RUN+="/home/tom/.local/bin/udev-xmodmap.sh"
