# option/mydns/mydns.service
[Unit]
Description=Notify IP to MyDNS

[Service]
Type=simple
ExecStart=/usr/bin/curl -s http://(MyDNSのID):(MyDNSのPass)@www.mydns.jp/login.html